# crontab
# Run the ETL script every 6 hours
0 */6 * * * cd /app && /usr/local/bin/python etl.py >> /var/log/cron.log 2>&1
# Run the ETL script at boot
@reboot cd /app && /usr/local/bin/python etl.py >> /var/log/cron.log 2>&1
